# deployment/reminder.conf
# Goes in /etc/init/

description reminder

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid www-data
setgid nogroup
chdir {{project_root}}/code/

exec {{project_root}}/env/bin/gunicorn {{app_name}}.wsgi --log-file=/tmp/{{app_name}}.log --log-level=debug --timeout=300